# .config/autostart/waypaper.desktop
[Desktop Entry]
Name=Waypaper
Comment=Wallpaper app
Exec=toolbox run --container waypaper waypaper
Type=Application
